# Andross/crontab
SHELL=/bin/bash
*/20 * * * * /usr/local/bin/python3.10 /Andross/create_snapshot.py -d > /dev/stdout
0 59 * * * /usr/local/bin/python3.10 /Andross/create_snapshot.py -l > /dev/stdout
